# root/etc/services.d/bazarr/run
#!/bin/sh
# Bazarr s6 service (latest version)

export HOME=/config

cd /app
exec /usr/local/bin/s6-setuidgid bsd python3.11 bazarr.py --no-update --config /config
